# supervisord.conf
[program:logstash-forwarder]
command = /opt/logstash-forwarder/logstash-forwarder -config /etc/forwarder/forwarder.conf
autorestart = true